# CamPanel.service
[Unit]
Description=CamPanel
After=network.target

[Service]
WorkingDirectory=/home/campanel/CamPanel2/
ExecStart=/usr/bin/python /home/campanel/CamPanel2/CamPanel.py
Restart=always
RestartSec=10
User=campanel

[Install]
WantedBy=multi-user.target